# weather-monitor.service
[Unit]
Description=Weather Monitor Service - SenseHAT Data Collector
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/weather
ExecStart=/usr/bin/python3 /home/pi/weather/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional)
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
